# deploy/bot-1c.ru/portal-api.service
# Скопировать: sudo cp portal-api.service /etc/systemd/system/portal-api.service
# Правки: User, Group, WorkingDirectory, EnvironmentFile

[Unit]
Description=Portal MVP FastAPI (uvicorn)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/www-root/data/www/bot-1c-api
EnvironmentFile=-/etc/default/portal-api
ExecStart=/var/www/www-root/data/www/bot-1c-api/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
